# files/etc/uci-defaults/00-zram-tuning
#!/bin/sh
# ZRAM Performance Tuning for 360V6 (IPQ6018)
LOG="/tmp/zram-tuning.log"
exec >> "$LOG" 2>&1
echo "=== ZRAM tuning start: $(date) ==="

# 1. 加载模块
if [ ! -d /sys/block/zram0 ]; then
    if ! modprobe zram num_devices=1; then
        echo "Error: zram module load failed"
        exit 1
    fi
    sleep 1
fi

# 2. 幂等检查
if swapon --show 2>/dev/null | grep -q zram0; then
    echo "ZRAM already active, skipping"
    exit 0
fi

# 3. 动态计算磁盘大小（物理内存的 50%，上限 512M）
TOTAL_MEM_KB=$(grep MemTotal /proc/meminfo | awk '{print $2}')
ZRAM_MB=$(( TOTAL_MEM_KB / 1024 / 2 ))
[ "$ZRAM_MB" -gt 512 ] && ZRAM_MB=512
[ "$ZRAM_MB" -lt 64  ] && ZRAM_MB=64
echo "Memory: ${TOTAL_MEM_KB}KB  ZRAM target: ${ZRAM_MB}MB"

# 4. 选择最佳压缩算法（优先 zstd > lz4 > lzo）
ALGO="lzo"
if [ -f /sys/block/zram0/comp_algorithm ]; then
    AVAIL=$(cat /sys/block/zram0/comp_algorithm)
    case "$AVAIL" in
        *zstd*) ALGO="zstd" ;;
        *lz4*)  ALGO="lz4"  ;;
    esac
fi
echo "$ALGO" > /sys/block/zram0/comp_algorithm \
    || { echo "Error: cannot set algorithm $ALGO"; exit 1; }
echo "Algorithm: $ALGO"

# 5. 设置大小
echo "${ZRAM_MB}M" > /sys/block/zram0/disksize \
    || { echo "Error: cannot set disksize"; exit 1; }

# 6. 格式化 + 激活
mkswap /dev/zram0 \
    || { echo "Error: mkswap failed"; exit 1; }
swapon /dev/zram0 -p 5 \
    || { echo "Error: swapon failed"; exit 1; }

echo "ZRAM enabled: ${ZRAM_MB}MB  algo=$ALGO"
exit 0
