# files/etc/uci-defaults/01-network-fix
#!/bin/sh
# 360V6 DSA Network Port Mapping Fix (严谨版)
# 硬件实际情况：1 WAN + 3 LAN，共4口
# 系统接口名：wan, lan1, lan2, lan3

LOG="/tmp/network_fix.log"
exec >> "$LOG" 2>&1
echo "=== Network fix start: $(date) ==="

# 1. 等待设备就绪
sleep 5

# 2. 强制重置网络配置
echo "Resetting network config for 360V6 (1WAN + 3LAN)..."

# 清除旧配置
uci delete network.lan
uci delete network.wan
uci delete network.device

# --- 配置 LAN (桥接 lan1, lan2, lan3) ---
uci set network.device=device
uci set network.device.name='br-lan'
uci set network.device.type='bridge'
# 修正：只添加存在的 lan1, lan2, lan3
uci add_list network.device.ports='lan1'
uci add_list network.device.ports='lan2'
uci add_list network.device.ports='lan3'

uci set network.lan=interface
uci set network.lan.device='br-lan'
uci set network.lan.proto='static'
# 修正 IP 避免冲突
uci set network.lan.ipaddr='192.168.2.1'
uci set network.lan.netmask='255.255.255.0'

# --- 配置 WAN (使用物理 wan 口) ---
uci set network.wan=interface
uci set network.wan.device='wan'
uci set network.wan.proto='dhcp'

# 3. 提交更改
uci commit network
echo "Network config committed."

# 4. 修复防火墙 (确保 wan 接口在 wan zone)
uci del_list firewall.@zone[1].network='wan' 2>/dev/null
uci add_list firewall.@zone[1].network='wan'
uci commit firewall

# 5. 重启网络生效
/etc/init.d/network restart
echo "=== Network fix done ==="
exit 0
