# files/etc/uci-defaults/01-network-fix
#!/bin/sh
# 360V6 DSA Network Port Mapping Fix

LOG="/tmp/network_fix.log"
exec >> "$LOG" 2>&1
echo "=== Network fix start: $(date) ==="

# 1. 等待 DSA 端口就绪
READY=false
for i in $(seq 1 30); do
    [ -d /sys/class/net/lan1 ] && { READY=true; break; }
    sleep 1
done

if [ "$READY" = "false" ]; then
    echo "Warning: LAN ports not detected"
fi

# 2. 动态发现端口
LAN_PORTS=""
for p in lan1 lan2 lan3 lan4; do
    [ -d "/sys/class/net/$p" ] && LAN_PORTS="$LAN_PORTS $p"
done
[ -z "$LAN_PORTS" ] && LAN_PORTS="lan1 lan2 lan3 lan4"

WAN_PORT=""
for p in wan eth4; do
    [ -d "/sys/class/net/$p" ] && WAN_PORT="$p" && break
done
[ -z "$WAN_PORT" ] && WAN_PORT="wan"

# 3. 配置 Network
uci -q batch << EOF
delete network.device
set network.device=device
set network.device.name=br-lan
set network.device.type=bridge
set network.lan=interface
set network.lan.device=br-lan
set network.lan.proto=static
set network.lan.ipaddr=192.168.2.1
set network.lan.netmask=255.255.255.0
set network.wan=interface
set network.wan.proto=dhcp
set network.wan.device=${WAN_PORT}
EOF

for p in $LAN_PORTS; do
    uci add_list network.device.ports="$p"
done

uci commit network
echo "network config committed"

# 4. 配置 Firewall WAN Zone
uci -q del_list firewall.@zone[1].network='wan'
uci add_list firewall.@zone[1].network='wan'
uci commit firewall

/etc/init.d/network reload 2>/dev/null || true
echo "=== Network fix done ==="
exit 0
