# files/etc/uci-defaults/99-custom-setup
#!/bin/sh

echo "=== Running 99-custom-setup ===" >> /tmp/setup.log
date >> /tmp/setup.log

# 确保配置文件存在
[ -f /etc/config/network ] || { echo "No network config!" >> /tmp/setup.log; exit 0; }

# 等待 DSA 初始化完成（关键！）
sleep 5

# 检查是否存在 lan 接口，如果不存在则创建（兼容 Master 分支自动生成逻辑）
if ! uci -q get network.lan > /dev/null 2>&1; then
    echo "Creating lan interface..." >> /tmp/setup.log
    uci set network.lan=interface
    uci set network.lan.proto='static'
    # DSA 下通常是 br-lan
    uci set network.lan.device='br-lan'
fi

# 确保 br-lan 桥接存在且包含物理端口
if ! uci -q get network.@device[0] > /dev/null 2>&1; then
    echo "Creating br-lan bridge..." >> /tmp/setup.log
    uci add network device
    uci set network.@device[-1].name='br-lan'
    uci set network.@device[-1].type='bridge'
    uci add_list network.@device[-1].ports='lan1'
    uci add_list network.@device[-1].ports='lan2'
    uci add_list network.@device[-1].ports='lan3'
fi

# 设置 IP
uci set network.lan.ipaddr='192.168.2.1'
uci set network.lan.netmask='255.255.255.0'

# 可选：禁用 DHCP（如果做旁路由/AP）
# uci set dhcp.lan.ignore='1'

uci commit network

echo "Setup completed" >> /tmp/setup.log

exit 0
