# uci-defaults/01_network_fix
#!/bin/sh
# ============================================================
# 360 V6 DSA Network Fix (无 Heredoc 版)
# ============================================================
LOG="/tmp/network_fix.log"
echo "Start Network Fix" > $LOG

# 1. 等待 DSA 端口就绪
for i in $(seq 1 30); do
    [ -d /sys/class/net/lan1 ] && break
    sleep 1
done

# 2. 动态收集端口 (严禁包含 eth0/eth1 CPU口)
LAN_PORTS=""
WAN_PORT=""
for p in lan1 lan2 lan3 lan4; do
    [ -d "/sys/class/net/$p" ] && LAN_PORTS="$LAN_PORTS $p"
done
[ -d "/sys/class/net/wan" ] && WAN_PORT="wan"

if [ -z "$LAN_PORTS" ]; then
    echo "Error: No LAN ports found" >> $LOG
    exit 1
fi

# 3. 写入网络配置 (使用 uci batch，避免 Heredoc)
uci batch << 'UCI_SCRIPT'
delete network.lan
delete network.wan
delete network.@device[0]
set network.@device[0]=device
set network.@device[0].name='br-lan'
set network.@device[0].type='bridge'
set network.lan=interface
set network.lan.device='br-lan'
set network.lan.proto='static'
set network.lan.ipaddr='192.168.2.1'
set network.lan.netmask='255.255.255.0'
set network.wan=interface
set network.wan.proto='dhcp'
UCI_SCRIPT

# 动态添加 LAN 端口到桥接
for p in $LAN_PORTS; do
    uci add_list network.@device[0].ports="$p"
done

# 设置 WAN 口
if [ -n "$WAN_PORT" ]; then
    uci set network.wan.device="$WAN_PORT"
fi

uci commit network
echo "Network config done" >> $LOG
exit 0
